# deploy/supervisor/roke-queue-worker.conf
; ─────────────────────────────────────────────────────────────────────────────
; Worker de colas de ROKE (Supervisor)
;
; Procesa las colas que mueven el aprovisionamiento y los despliegues. SIN este
; worker corriendo, las sagas de cómputo (RunOrchestration) y el aprovisionamiento
; de hosting (RunProvisioningJob) se quedan encolados para siempre → los servicios
; quedan "Creando…" sin avanzar (justo lo que pasó con angular-crud-frontend).
;
; Instalación (en el server donde corre Laravel — Mac Mini dev / Dell prod):
;   sudo cp deploy/supervisor/roke-queue-worker.conf /etc/supervisor/conf.d/
;   # Ajusta las rutas y el usuario abajo, luego:
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start roke-worker:*
;
; Tras cada despliegue de código:  php artisan queue:restart
; ─────────────────────────────────────────────────────────────────────────────

[program:roke-worker]
process_name=%(program_name)s_%(process_num)02d
; --queue: prioridad de izquierda a derecha. 'default' al final para el resto de jobs.
; --tries=1: los jobs (RunOrchestration / RunProvisioningJob) gestionan sus propios
;            reintentos/backoff; no queremos que Laravel los reintente por su cuenta.
; --max-time=3600: recicla el proceso cada hora para evitar fugas de memoria.
command=php /ruta/al/proyecto/artisan queue:work redis --queue=provisioning,deployments,default --sleep=3 --tries=1 --max-time=3600
directory=/ruta/al/proyecto
autostart=true
autorestart=true
stopwaitsecs=130
; Más alto que el timeout del job más largo (RunOrchestration::$timeout = 120s)
; para que un reinicio no mate un paso a medias.
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/ruta/al/proyecto/storage/logs/worker.log
stopasgroup=true
killasgroup=true
